# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:go2rtc]
command=/usr/local/bin/go2rtc -config /app/go2rtc.yaml
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/go2rtc.err.log
stdout_logfile=/var/log/go2rtc.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:config-api]
command=/usr/bin/python3 /app/api/config_api.py
directory=/app/api
autostart=true
autorestart=true
stderr_logfile=/var/log/config-api.err.log
stdout_logfile=/var/log/config-api.out.log

[program:status-api]
command=/usr/bin/python3 /app/api/status_api.py
directory=/app/api
autostart=true
autorestart=true
stderr_logfile=/var/log/status-api.err.log
stdout_logfile=/var/log/status-api.out.log
